# context/splynx-start.service
[Unit]
After=mysql.service

[Service]
ExecStart=/var/www/splynx/docker/start-systemd.sh

[Install]
WantedBy=default.target
